# dist/recmeet-server.socket
[Unit]
Description=Recmeet V2 Server Socket

[Socket]
ListenStream=%t/recmeet-server/server.sock
SocketMode=0700

[Install]
WantedBy=sockets.target
