# dist/recmeet-daemon.socket
[Unit]
Description=Recmeet Daemon Socket

[Socket]
ListenStream=%t/recmeet/daemon.sock
SocketMode=0700

[Install]
WantedBy=sockets.target
