# dist/recmeet-daemon.service.in
[Unit]
Description=Recmeet Meeting Recorder Daemon
Documentation=https://github.com/suykerbuyk/recmeet

[Service]
Type=simple
ExecStart=@RECMEET_DAEMON_EXEC@
ExecReload=/bin/kill -HUP $MAINPID
Environment=RECMEET_LOG_LEVEL=error
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
