# dist/recmeet-client.service.in
[Unit]
Description=Recmeet V2 Client (system tray / WebUI bridge / audio capture)
Documentation=https://github.com/suykerbuyk/recmeet
PartOf=graphical-session.target
Wants=recmeet-server.service
After=graphical-session.target recmeet-server.service

[Service]
Type=simple
ExecStart=@RECMEET_CLIENT_EXEC@
Environment=RECMEET_LOG_LEVEL=error
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
